# crontab/crontab.docker
00 13 * * 1 docker rmi $(docker images | grep none | awk '{print $3}')
00 13 * * 1 docker volume rm $(docker volume ls -qf dangling=true)
