# crontab/install.sh
#!/bin/bash
set -e
source "$DOTFILES/.lib/include.sh"
actualize_crontabs
